# systemd/caddy.service.d/cloudflare.conf
# Wires the Cloudflare API token into caddy.service environment so the
# wildcard *.sacred.vote DNS-01 ACME challenge can call the Cloudflare API.
# The env file at /etc/caddy/cloudflare.env is mode 0640 caddy:caddy.
# This drop-in is part of #256 (wildcard TLS cert) and is staged but
# inert until the Caddyfile imports `wildcard_tls` inside its sites.
[Service]
EnvironmentFile=/etc/caddy/cloudflare.env
