# systemd/caddy.service.d/onfailure.conf
[Unit]
# Page the operator if Caddy dies — this is the public front door for
# every site (sacred.vote, sacredvote.org, mail.sacred.vote,
# vault.sacred.vote, files.sacred.vote). Total outage if it dies.
OnFailure=failure-alert@%n.service
